# api-server/etc/services.d/vpn-api/run
#!/usr/bin/with-contenv bash

# Wait for NordVPN to be ready
sleep 15

# Start the Python VPN API server as root (needed for nordvpn commands)
exec python3 /usr/bin/vpn_api_server_py 8080